# ops/systemd/amonora-core-pg-backup.service
[Unit]
Description=Amonora core PostgreSQL backup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/amonora/backup.env
ExecStart=/bin/bash /opt/amonora_bot/ops/backup/core_pg_backup.sh

[Install]
WantedBy=multi-user.target
