# ops/systemd/amonora-grafana.service
[Unit]
Description=Amonora Grafana OSS analytics
After=network-online.target amonora-grafana-db-tunnel.service
Wants=network-online.target
Requires=amonora-grafana-db-tunnel.service

[Service]
Type=simple
User=grafana
Group=grafana
WorkingDirectory=/opt/amonora_bot
EnvironmentFile=-/etc/amonora/grafana.env
ExecStart=/usr/sbin/grafana-server \
  --config=/opt/amonora_bot/ops/grafana/grafana.ini \
  --homepath=/usr/share/grafana \
  cfg:default.paths.data=/var/lib/grafana \
  cfg:default.paths.logs=/var/log/grafana \
  cfg:default.paths.plugins=/var/lib/grafana/plugins \
  cfg:default.paths.provisioning=/opt/amonora_bot/ops/grafana/provisioning
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
MemoryMax=512M
CPUQuota=30%

[Install]
WantedBy=multi-user.target
