# ops/systemd/amonora-grafana-db-tunnel.service
[Unit]
Description=Amonora Grafana PostgreSQL tunnel to core
After=network-online.target
Wants=network-online.target
Before=amonora-grafana.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/amonora_bot
EnvironmentFile=-/etc/amonora/grafana-db-tunnel.env
ExecStart=/bin/bash -lc 'exec /usr/bin/ssh -N \
  -L "${GRAFANA_ANALYTICS_TUNNEL_BIND_HOST:-127.0.0.1}:${GRAFANA_ANALYTICS_TUNNEL_LOCAL_PORT:-15432}:127.0.0.1:5432" \
  -i "${GRAFANA_ANALYTICS_TUNNEL_KEY_PATH}" \
  -o BatchMode=yes \
  -o ExitOnForwardFailure=yes \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o StrictHostKeyChecking=yes \
  -o UserKnownHostsFile="${GRAFANA_ANALYTICS_TUNNEL_KNOWN_HOSTS}" \
  -o ConnectTimeout="${GRAFANA_ANALYTICS_TUNNEL_CONNECT_TIMEOUT:-15}" \
  -p "${GRAFANA_ANALYTICS_TUNNEL_CORE_PORT:-22}" \
  "${GRAFANA_ANALYTICS_TUNNEL_CORE_USER}@${GRAFANA_ANALYTICS_TUNNEL_CORE_HOST}"'
Restart=always
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
